# t-rest.service
[Unit]
Description=Tyrannasaurus Rest

[Service]
Type=simple
ExecStart=/usr/local/bin/deno run --allow-run 

[Install]
WantedBy=multi-user.target
